# script/docker/xray.ini
[program:xray]
command=xray run -config /usr/local/etc/xray/config.json
autostart=true
autorestart=unexpected
startsecs=1
directory=/tmp
stdout_logfile=/tmp/xray.log
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 10
